# deploy/systemd/btc-flows.timer
[Unit]
Description=Schedule Farside flow refresh

[Timer]
OnCalendar=*-*-* 23,01,02:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
