# deploy/systemd/btc-flows.timer
[Unit]
Description=Schedule Farside flow refresh

[Timer]
OnCalendar=*-*-* 22,23,01:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
