# deploy/systemd/farside-flows.timer
[Unit]
Description=Schedule Farside flow refresh

[Timer]
OnCalendar=*-*-* 03,04,05:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
